# riemann/init/riemann-dash
#!/sbin/runscript

# This script assumes that Riemann has a user and group set up, and that the
# system gem `riemann-dash` has been installed.
#
# Author: Greg Look (greg@mvxcvi.com)

CFG_DIR=/etc/riemann
RUN_DIR=/var/run/riemann
LOG_DIR=/var/log/riemann

COMMAND=/usr/local/bin/riemann-dash

CONFIG_FILE=$CFG_DIR/dash.rb
PIDFILE=$RUN_DIR/dash.pid
LOGFILE=$LOG_DIR/dash.log

depend() {
    use riemann net
}

start() {
    checkpath -d -m 0755 -o riemann:riemann $RUN_DIR
    checkpath -d -m 0755 -o riemann:riemann $LOG_DIR
    ebegin "Starting riemann dashboard"
    start-stop-daemon \
        --start \
        --user riemann:riemann \
        --chdir $RUN_DIR \
        --pidfile $PIDFILE \
        --make-pidfile \
        --background \
        --stdout $LOGFILE \
        --stderr $LOGFILE \
        --exec $COMMAND \
        -- $CONFIG_FILE
    eend $?
}

stop() {
    ebegin "Stopping riemann dashboard"
    start-stop-daemon \
        --stop \
        --pidfile $PIDFILE \
        --exec $COMMAND
    eend $?
}
